# debian/ifupdown2.networking.service
[Unit]
Description=ifupdown2 init script

[Service]
Type=oneshot
ExecStart=/sbin/ifup -a
ExecStop=/sbin/ifdown -a

[Install]
WantedBy=multi-user.target
